# debian/postinst
#!/bin/bash
# postinst script for myenv
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

MYSQL_GROUP=mysql
MYSQL_USER=mysql
MYSQL_HOME=/home/mysql

function configure_files ()
{
cd ${MYSQL_HOME}/product
rm -f myenv                   # On upgrade, old files are already deleted, but myenv-RELEASE/log still exists
ln -s myenv-2.0.0 myenv   # Remove symlink, then set it again.

# Make the "myenv" functions accessible in the user's shell
ed ${MYSQL_HOME}/.bash_profile <<EOF 2>/dev/null || :
/^# BEGIN MyEnv$/,/^# END MyEnv$/d
w
q
EOF
cat myenv/tpl/profile.template >> ${MYSQL_HOME}/.bash_profile

echo "export MYENV_BASE=${MYSQL_HOME}/product/myenv" > /etc/myenv/MYENV_BASE
chown -R ${MYSQL_USER}:${MYSQL_GROUP}  \
    ${MYSQL_HOME}/product/ ${MYSQL_HOME}/.bash_profile /etc/myenv
}

function handle_startup ()
{
if [ -d /run/systemd/system ] ; then
    install -D -m 0644 ${MYSQL_HOME}/product/myenv/tpl/systemd.myenv.unit.template /etc/systemd/system/myenv.service
    systemctl daemon-reload >/dev/null 2>&1 || :
    systemctl enable myenv >/dev/null 2>&1 || :
else
    install -D -m 0755 ${MYSQL_HOME}/product/myenv/tpl/myenv.server /etc/init.d/myenv
    # Will be activated by the debhelper snippet at the end
fi
}

case "$1" in
    configure)
        configure_files
        handle_startup
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        : # no action
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
        if [ -x "/etc/init.d/myenv" ]; then
                update-rc.d myenv defaults >/dev/null
        fi
#       if [ -x "/etc/init.d/myenv" ] || [ -e "/etc/init/myenv.conf" ]; then
#               invoke-rc.d myenv start || exit $?
#       fi
fi
# End automatically added section


exit 0
